# distros/raspberry_pi/lighthouse/w3p_lighthouse-beacon.service
[Unit]
Description=Web3Pi Lighthouse Consensus Layer Client service
After=network.target
Documentation=https://setup-guide.web3pi.io

[Service]
Type=simple
WorkingDirectory=/home/ethereum
ExecStart=/home/ethereum/clients/lighthouse/lighthouse.sh
User=ethereum
Restart=on-failure
RestartSec=120
KillSignal=SIGINT
TimeoutStopSec=900

[Install]
WantedBy=multi-user.target
